# systemd/kaguya-mcp.service
[Unit]
Description=Kaguya MemPalace MCP Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/apps/kaguya-gateway
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/home/ubuntu
Environment=PATH=/home/ubuntu/apps/kaguya-gateway/.venv/bin:/home/ubuntu/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/home/ubuntu/apps/kaguya-gateway/.env
ExecStart=/home/ubuntu/apps/kaguya-gateway/.venv/bin/python -m app.mcp.server
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
